# downloader/etc/supervisord.conf
[unix_http_server]
file=%(here)s/../var/run/supervisor.sock   ; path to your socket file
;chmod=0777
;chown=webscan:webscan

[supervisord]
logfile=%(here)s/../var/log/supervisord.log ; supervisord log file
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=info               ; info, debug, warn, trace
pidfile=%(here)s/supervisord.pid ; pidfile location
nodaemon=false              ; run supervisord as a daemon
minfds=1024                 ; number of startup file descriptors
minprocs=200                ; number of process descriptors
user=root                  ; default user
childlogdir=%(here)s          ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/../var/run/supervisor.sock
prompt=supervisor/downloader
history_file=%(here)s/../var/cache/supervisorctl.history

[program:downloader]
command=python %(here)s/../bin/downloader.py -f %(here)s/downloader.conf
autorestart=true
autostart=true
startsecs=1
startretries=10000
stopwaitsecs=3
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=syslog
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
directory=%(here)s/../
environment=ST_AUTH=http://10.171.55.16:8080/auth/v1.0,ST_USER=system:root,ST_KEY=testpass
